# tests/supervisord/celery.conf
[program:celery-training]
command=/Users/johndoe/anaconda3/envs/linc-cv/bin/celery worker -A linc_cv.tasks --concurrency 1 -Q training -E -n training@%%h
environment=PYTHONPATH=/Users/johndoe/work/git/linc/linc-cv,API_KEY=blah
user=johndoe
redirect_stderr=true
stdout_logfile=/Users/johndoe/work/git/linc/linc-cv/tests/celery-training.log
stopasgroup=true
killasgroup=true

[program:celery-classification]
command=/Users/johndoe/anaconda3/envs/linc-cv/bin/celery worker -A linc_cv.tasks --concurrency 1 -Q classification -E -n classification@%%h
environment=PYTHONPATH=/Users/johndoe/work/git/linc/linc-cv,API_KEY=blah
user=johndoe
redirect_stderr=true
stdout_logfile=/Users/johndoe/work/git/linc/linc-cv/tests/celery-classification.log
stopasgroup=true
killasgroup=true
