# tests/supervisord/celery.conf
[program:celery-training]
command=/Users/sglee/anaconda3/envs/linc-cv/bin/celery worker -A linc_cv.tasks --concurrency 1
environment=PYTHONPATH=/Users/sglee/work/git/linc/linc-cv,API_KEY=blah
user=sglee
redirect_stderr=true
stdout_logfile=/Users/sglee/work/git/linc/linc-cv/tests/celery-training.log
stopasgroup=true
killasgroup=true

[program:celery-classification]
command=/Users/sglee/anaconda3/envs/linc-cv/bin/celery worker -A linc_cv.tasks --concurrency 1
environment=PYTHONPATH=/Users/sglee/work/git/linc/linc-cv,API_KEY=blah
user=sglee
redirect_stderr=true
stdout_logfile=/Users/sglee/work/git/linc/linc-cv/tests/celery-classification.log
stopasgroup=true
killasgroup=true
