# tests/supervisord/flower.conf
[program:flower]
command=/Users/sglee/anaconda3/envs/linc-cv/bin/flower -A linc_cv.tasks --address=0.0.0.0 --port=5555
environment=PYTHONPATH=/Users/sglee/work/git/linc/linc-cv,API_KEY=blah
user=sglee
redirect_stderr=true
stdout_logfile=/Users/sglee/work/git/linc/linc-cv/tests/flower.log
stopasgroup=true
killasgroup=true