# tests/supervisord/flower.conf
[program:flower]
command=/Users/johndoe/anaconda3/envs/linc-cv/bin/celery -A linc_cv.tasks flower --address=0.0.0.0 --port=5555
environment=PYTHONPATH=/Users/johndoe/work/git/linc/linc-cv,API_KEY=blah
user=johndoe
redirect_stderr=true
stdout_logfile=/Users/johndoe/work/git/linc/linc-cv/tests/flower.log
stopasgroup=true
killasgroup=true